# examples/elbe_advanced/image/rpi-image/overlays/systemd/etc/systemd/system/coffee.service
[Unit]
Description=Coffee service

[Service]
ExecStart=/usr/examples/coffee

[Install]
WantedBy=multi-user.target
