# examples/yocto_advance/qt/meta-coffee/recipes-qt/qt6/qtdoc/coffee.service
[Unit]
Description=Coffeemachine

[Service]
ExecStart=/usr/share/examples/demos/coffee/coffeemachine

[Install]
WantedBy=multi-user.target
